# Utils/StarPipeline/SystemDTemplates/populate-work-queue.timer
[Unit]
Description=Run PopulateWorkQueue every hour

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
