# Utils/StarPipeline/SystemDTemplates/sort-calstar-cache.timer
[Unit]
Description=Run CALSTAR cache sort every hour

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target

