# Utils/StarPipeline/SystemDTemplates/star-pipeline.service
[Unit]
Description=StarPipeline Ingestion Service
After=network-online.target
Wants=network-online.target

[Service]
User=star
WorkingDirectory=%h/source/RMS
ExecStart=%h/source/RMS/Scripts/RunStarPipeline.sh
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
